# scripts/systemd/ai-nids-dashboard.service
[Unit]
Description=AI-NIDS Live Dashboard
After=network.target ai-nids-monitor.service

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{PROJECT_ROOT}}

# Explicitly bind Streamlit to correct host/port if necessary
ExecStart="{{PROJECT_ROOT}}/ai-venv/bin/streamlit" run dashboard/app.py --server.port 8501 --server.address 127.0.0.1
Restart=always
RestartSec=5

Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
