# scripts/systemd/ai-nids-monitor.service
[Unit]
Description=AI-NIDS Inspection Engine Driver
After=network.target

[Service]
Type=simple
User=root
# Replace /home/dev/Projects/ai_nids with your actual deployment path if moved!
WorkingDirectory=/home/dev/Projects/ai_nids

# Assuming the virtual environment is inside the project dir: 'ai-venv'
ExecStart=/home/dev/Projects/ai_nids/ai-venv/bin/python scripts/run_monitor.py --interface eth0 --flow-timeout 20
Restart=always
RestartSec=5

# Helps Python logs stream to journalctl immediately
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
