# scripts/systemd/ai-nids-dashboard.service
[Unit]
Description=AI-NIDS Live Dashboard
After=network.target ai-nids-monitor.service

[Service]
Type=simple
User=dev
WorkingDirectory=/home/dev/Projects/ai_nids

# Explicitly bind Streamlit to correct host/port if necessary
ExecStart=/home/dev/Projects/ai_nids/ai-venv/bin/streamlit run dashboard/app.py --server.port 8501 --server.address 127.0.0.1
Restart=always
RestartSec=5

Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
